# server/Docker_twisted/supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=twistd -r epoll web --class=pyapns.server.APNSServer --port=80
redirect_stderr=true
autostart=true
autorestart=true

