# server/Docker_ippuq_api/supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=python ippuq_api.py
redirect_stderr=true
autostart=true
autorestart=true

